# dist/submarine.service
[Unit]
Description=Submarine
Wants=network-online.target
After=network-online.target

[Service]
User=pi
Group=pi
Type=simple
Restart=on-failure
RestartSec=3
WorkingDirectory=/home/pi/submarine
AmbientCapabilities=CAP_SYS_NICE
Environment="RUST_LOG=info"
ExecStart=/home/pi/submarine/submarine

[Install]
WantedBy=multi-user.target